# etc/systemd/system/proxywise.service
[Unit]
Description=Proxywise

[Service]
ExecStart=/usr/bin/proxywise-service

[Install]
WantedBy=multi-user.target

