# etc/systemd/system/proxywise.service
[Unit]
Description=Proxywise

[Service]
ExecStart=proxywise-service

[Install]
WantedBy=multi-user.target

